# build/supervisord.conf
[supervisord]
nodaemon=true

[program:validator]
command=validator --beacon-rpc-provider=prysm-beacon-chain.public.dappnode:4000 --keystore-path=/data --password=%(ENV_PASSWORD)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0